# src/frontend/diagon_alley.socket
[Unit]
Description=FAUST CTF diagon alley service

[Socket]
ListenStream=4441
Accept=true

[Install]
WantedBy=sockets.target
